# admin/userdb.socket
Unit]
Description=gunicorn socket

[Socket]
ListenStream=/var/run/userdb.sock

[Install]
WantedBy=sockets.target
